# setup-guide/docker.service.d/override.conf
[Unit]
After=rclone-media.service mnt-ebs\x2dmp.mount
Requires=rclone-media.service mnt-ebs\x2dmp.mount
